# no-kmod-autoload.service.in
# SPDX-FileCopyrightText: 2026 Julian Brost <julian@0x4a42.net>
# SPDX-License-Identifier: MIT

[Unit]
Description=Disable automatic module loading after boot finished

[Service]
# systemd waits for oneshot, use simple so that the unit isn't waiting for itself
Type=simple
ExecStart=@libexecdir@/no-kmod-autoload/start.sh @libexecdir@/no-kmod-autoload/modprobe.sh

[Install]
WantedBy=multi-user.target
